# outputs/oci/bik-toss-collector-eod.service
[Unit]
Description=BIK Research Toss EOD candles collector
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/bik-research/outputs
EnvironmentFile=/etc/bik-toss-collector.env
ExecStart=/usr/bin/env TOSS_COLLECT_KR_PRICES=true TOSS_COLLECT_KR_STOCK_INFO=false TOSS_COLLECT_KR_PRICE_LIMITS=true TOSS_COLLECT_KR_CANDLES=true TOSS_CANDLE_INTERVALS=1d TOSS_CANDLE_COUNT=60 TOSS_DETAIL_SYMBOL_SCOPE=all TOSS_DETAIL_SYMBOL_LIMIT=0 TOSS_DETAIL_INTERVAL_SECONDS=0 TOSS_BATCH_SLEEP_SECONDS=0.2 /home/ubuntu/bik-research/outputs/.venv/bin/python /home/ubuntu/bik-research/outputs/toss_collector.py
